# cdk/cdk/scripts/game-server.service
# Defines the game WebSocket server as a systemd service.

# This file should be copied to /etc/systemd/system/ during bootstrap,
# then enabled with: systemctl enable game-server.service
# so it starts automatically on every boot.

# When started, the service authenticates with ECR, pulls the latest image,
# and runs it as a Docker container, restarting automatically if it crashes.

[Unit]
Description=Game WebSocket Server
# Start this service after the network is fully up and Docker is running
After=network-online.target docker.service update-dns.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
# Authenticate with ECR before pulling
# See https://docs.aws.amazon.com/AmazonECR/latest/userguide/registry_auth.html
ExecStartPre=/bin/bash -c 'aws ecr get-login-password --region {{AWS_REGION}} | docker login --username AWS --password-stdin {{ECR_REGISTRY}}'
# Pull the latest image from ECR
ExecStartPre=/usr/bin/docker pull {{ECR_REPO_URI}}:latest
# Run the container
ExecStart=/usr/bin/docker run --rm --name game-server --network game-network -p {{GAME_SERVER_PORT}}:{{GAME_SERVER_PORT}} {{ECR_REPO_URI}}:latest
# Automatically recover if the server crashes, but leave it stopped if it shut down intentionally
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target